# service/script/up
#!/usr/bin/env bash

function main {
  local shome="$(cd -P -- "${BASH_SOURCE%/*}/.." && pwd -P)"

  touch "$shome/.connected"
}

main "$@"
